# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

; --- infrastructure -------------------------------------------------------
[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql
user=root
priority=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --save "" --appendonly no
user=redis
priority=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- one-shot site bootstrap ---------------------------------------------
[program:site-init]
command=/usr/local/bin/site-init.sh
user=root
priority=5
autostart=true
autorestart=false
startsecs=0
startretries=1
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Frappe / ERPNext application processes -------------------------------
[program:backend]
command=/bin/bash -c 'cd /home/frappe/frappe-bench && exec env/bin/gunicorn --chdir=/home/frappe/frappe-bench/sites --bind=127.0.0.1:8000 --workers=%(ENV_GUNICORN_WORKERS)s --threads=4 --worker-class=gthread --worker-tmp-dir=/dev/shm --timeout=120 frappe.app:application'
user=frappe
environment=HOME="/home/frappe",USER="frappe",PATH="/home/frappe/frappe-bench/env/bin:/home/frappe/.local/bin:/usr/local/bin:/usr/bin:/bin"
priority=10
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:socketio]
command=/bin/bash -c 'cd /home/frappe/frappe-bench && exec node apps/frappe/socketio.js'
user=frappe
environment=HOME="/home/frappe",USER="frappe",PATH="/home/frappe/frappe-bench/env/bin:/home/frappe/.local/bin:/usr/local/bin:/usr/bin:/bin"
priority=10
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=/bin/bash -c 'cd /home/frappe/frappe-bench && exec bench schedule'
user=frappe
environment=HOME="/home/frappe",USER="frappe",PATH="/home/frappe/frappe-bench/env/bin:/home/frappe/.local/bin:/usr/local/bin:/usr/bin:/bin"
priority=15
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/bin/bash -c 'cd /home/frappe/frappe-bench && exec bench worker'
user=frappe
environment=HOME="/home/frappe",USER="frappe",PATH="/home/frappe/frappe-bench/env/bin:/home/frappe/.local/bin:/usr/local/bin:/usr/bin:/bin"
priority=15
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- web frontend ---------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
